# infra/oci-gallery-agent/already-here-drive-sync.timer
[Unit]
Description=Run Already Here approved gallery Drive sync every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true
Unit=already-here-drive-sync.service

[Install]
WantedBy=timers.target
